# systemd/hermes-content-platform-maintenance.timer
[Unit]
Description=Run Hermes content maintenance daily

[Timer]
OnCalendar=*-*-* 04:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

