# systemd/hermes-content-platform-task-market.timer
[Unit]
Description=Run Hermes AiToEarn task-market automation every 4 hours

[Timer]
OnCalendar=*-*-* 00/4:15:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
